# daemon/systemd/mud.service
[Unit]
Description=Mu Daemon Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mud

[Install]
WantedBy=default.target
